# openrc-net-namespaces/src/etc/init.d/net_namespace
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need localmount net
	use dns
	after bootmisc
}

start() {
	ebegin "Configuring $NS namespace"
	ip netns add $NS
	ip link add $HOST_VETH type veth peer name $NS_VETH
	ip link set $NS_VETH netns $NS
	brctl addif $HOST_BRIDGE $HOST_VETH
	ifconfig $HOST_VETH up
	ip netns exec $NS ifconfig $NS_VETH up
	ip netns exec $NS ifconfig lo up
	ip netns exec $NS ip addr add $NS_IP dev $NS_VETH
	#forwarding
	ip netns exec $NS sysctl net.ipv4.ip_forward=1
	for route in $NS_ROUTES; do
		ip netns exec $NS ip route add $route via $NS_ROUTES_GW
	done
	eend 0
	for service in $NS_SERVICES; do
		ip netns exec $NS /etc/init.d/$service start
        done
	return 0
}

stop() {
	for service in $NS_SERVICES; do
		ip netns exec $NS /etc/init.d/$service stop
	done
	
	ebegin "Deleting $NS namespace"
	ip link del $HOST_VETH
	ip netns delete $NS
	eend 0
}

# vim: set ts=4 :
